# simple-c2/crontab
5 */1 * * * /home/handyman/simple-c2/c2-client.sh