# ip-bacon/crontab
0 */2 * * * /home/handyman/ip-bacon/ip-bacon.sh