# dns-stats/crontab
* * * * * /home/handyman/dns-stats/check-dns.sh